# systemd/.config/systemd/user/muindex.timer
[Unit]
Description=mu index timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=60m
Unit=muindex.service

[Install]
WantedBy=timers.target
